# Raspi5/services/homeeye-web.service
[Unit]
Description=HomeEye web server
After=network.target

[Service]
User=homeeye
EnvironmentFile=/opt/homeeye/Raspi5/.env
WorkingDirectory=/opt/homeeye/server
ExecStart=/usr/bin/node /opt/homeeye/server/server.js
Restart=always
RestartSec=5
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
[Unit]
Description=HomeEye Web UI + API (Express)
After=network.target


[Service]
Type=simple
User=riplay
EnvironmentFile=/home/riplay/homeeye/Raspi5/.env
WorkingDirectory=/home/riplay/homeeye/server
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=2


[Install]
WantedBy=multi-user.target
