# Raspi5/services/homeeye-camera.service
[Unit]
Description=HomeEye camera streamer
After=network.target

[Service]
User=homeeye
EnvironmentFile=/opt/homeeye/Raspi5/.env
WorkingDirectory=/opt/homeeye/Raspi5/camera
ExecStart=/usr/bin/node /opt/homeeye/Raspi5/camera/streamer.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
[Unit]
Description=HomeEye Camera MJPEG Streamer (libcamera-vid → Node)
After=network.target


[Service]
Type=simple
User=riplay
EnvironmentFile=/home/riplay/homeeye/Raspi5/.env
WorkingDirectory=/home/riplay/homeeye
ExecStart=/usr/bin/node Raspi5/camera/streamer.js
Restart=on-failure
RestartSec=2


[Install]
WantedBy=multi-user.target
